# systemd/production/youtility-mqtt.service
[Unit]
Description=YOUTILITY MQTT Client (Production)
After=network.target
Wants=network.target

[Service]
Type=simple
User=redmine
Group=redmine
WorkingDirectory=/home/redmine/django5_project/YOUTILITY3
Environment=VIRTUAL_ENV=/home/redmine/django5_project/env
Environment=PATH=/home/redmine/django5_project/env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/redmine/django5_project/env/bin/python paho_client.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=youtility-mqtt

# Resource limits
LimitNOFILE=1024
MemoryHigh=256M
MemoryMax=512M

# Security settings (Production)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/django5

[Install]
WantedBy=multi-user.target